# deploy/jarvis-staging.service
[Unit]
Description=Jarvis Staging Agent
After=network.target
# systemd ignores these in [Service]; they only take effect in [Unit]. Give up
# after 5 failed starts within 5 min instead of restarting forever.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=jarvis_user
Group=jarvis_user
WorkingDirectory=/app/jarvis_staging/code
# Point directly to the Python executable inside the staging virtual environment
ExecStart=/app/jarvis_staging/code/venv/bin/python3 /app/jarvis_staging/code/main.py
Restart=always
RestartSec=30
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
# Ensure the staging venv is in the path
Environment="PATH=/app/jarvis_staging/code/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Instance identity only — no proactive toggles. All three default off in code, so
# staging is inert (no heartbeat, no reminders, no webhook) and never messages the
# owner. Add a single JARVIS_*_ENABLED=true line here only for a run that
# deliberately wants that behavior.
Environment="JARVIS_ROOT=/app/jarvis_staging"

[Install]
WantedBy=multi-user.target
